# misc/systemd/isg-civil-backup.timer
[Unit]
Description=Run Civil Database backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
