# misc/systemd/isg-civil-backup.timer
[Unit]
Description=Run Civil Database backup weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
