# misc/systemd/isg-civil-backup.timer
[Unit]
Description=Run Civil Database backup weekly

[Timer]
Unit=isg-civil-backup.service
OnCalendar=weekly

[Install]
WantedBy=timers.target
